# crontab
# Mondays and Thursdays at 08:00 — search for new events for all verified users
0 8 * * 1,3,5 cd /app && python run_pipeline.py
